# custom-os/airootfs/etc/profile.d/enable_sddm.sh
#!/bin/bash
# Enable SSDM to start on boot
systemctl enable sddm
