# custom-os/airootfs/etc/profile.d/install_hacknerd_font.sh
#!/bin/bash
# Run the Hack Nerd Font installation script
/root/scripts/install_hacknerd_font.sh

