# custom-os/airootfs/etc/profile.d/run_nessus_install.sh
#!/bin/bash
# Run Nessus installation script on live boot
if [ -f /root/install_nessus.sh ]; then
    /root/install_nessus.sh
fi
