# custom-os/airootfs/etc/profile.d/run_metasploit_setup.sh
#!/bin/bash
# Run Metasploit setup script on live boot
if [ -f /root/setup_metasploit.sh ]; then
    /root/setup_metasploit.sh
fi

