# custom-os/airootfs/etc/profile.d/run_dnsrecon_install.sh
#!/bin/bash
# Run DNSRecon installation script on live boot
if [ -f /root/install_dnsrecon.sh ]; then
    /root/install_dnsrecon.sh
fi

