# custom-os/airootfs/etc/profile.d/run_openvas_install.sh
#!/bin/bash
# Run OpenVAS installation script on live boot
if [ -f /root/install_openvas.sh ]; then
    /root/install_openvas.sh
fi
