# custom-os/airootfs/etc/profile.d/start_networkmanager.sh
#!/bin/bash
# Enable and start NetworkManager for networking in live ISO
systemctl enable NetworkManager
